# build/nobarrier.service
[Unit]
Description=Auto detect and connect to barrier server
After=network.target
StartLimitIntervalSec=0
[Service]
Type=simple
Restart=on-failure
RestartSec=1
User=pi
ExecStart=/home/pi/Documents/nobarrier/build/nobarrier pi

[Install]
WantedBy=multi-user.target
